# debian/ioclogserver.target
[Unit]
Description=EPICS IOC Log Server (multi-instance-master)

[Install]
WantedBy=multi-user.target
